# airootfs/etc/systemd/user/notifitor.service
[Service]
ExecStart=/usr/local/bin/notifitor

[Install]
WantedBy=graphical-session.target
