# airootfs/etc/systemd/system/plymouth-deactivate.service
[Unit]
Description=Deactivate Plymouth
After=getty.target

[Service]
ExecStart=/usr/bin/plymouth deactivate

[Install]
WantedBy=multi-user.target
