# airootfs/etc/systemd/system/switch-vt.service
[Unit]
Description=Switch to TTY2
After=getty.target

[Service]
ExecStart=/usr/bin/chvt 2

[Install]
WantedBy=multi-user.target
