# airootfs/etc/systemd/system/switch-vt.service
[Unit]
Description=Switch to TTY2
After=systemd-user-sessions.service

[Service]
ExecStart=/usr/bin/chvt 2

[Install]
WantedBy=multi-user.target
